# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:crond]
command=/usr/sbin/cron

[program:geoip]
command=/usr/bin/python3 /opt/geoip-flask/app.py
