# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:geoip]
command=/usr/local/bin/uwsgi --master --http 0.0.0.0:8888 --chdir /opt/geoip-flask/ --module app:app --processes 4 --stats 0.0.0.0:9191 --stats-http
